# deployment/svd_wsrtunnel.conf
[program:wsrtunnel]
process_name = wsrtunnel

command=python3 /app/server.py

user=nobody
autostart=true
autorestart=true
