# deployment/svd_nginx.conf
[program:nginx]
process_name = nginx

command=/usr/sbin/nginx

autostart=true
autorestart=true
